# storage-app.service
[Unit]
Description=Storage App - Skladová aplikace pro Gitterbox management
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/storage-app/backend
Environment=PATH=/opt/storage-app/venv/bin
Environment=PYTHONPATH=/opt/storage-app
ExecStart=/opt/storage-app/venv/bin/python start_production.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Bezpečnostní omezení
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/storage-app
ProtectHome=yes

# Logy
StandardOutput=journal
StandardError=journal
SyslogIdentifier=storage-app

[Install]
WantedBy=multi-user.target
